# config/files/usr/etc/profile.d/fzf.sh
shell=$(basename $(echo $SHELL))
source "/usr/share/fzf/shell/key-bindings.$shell"
