# config/files/usr/etc/profile.d/starship.sh
shell=$(basename $(echo $SHELL))
eval "$(starship init $(echo $shell))"
