# config/files/usr/etc/profile.d/zoxide.sh
shell=$(basename $(echo $SHELL))
eval "$(zoxide init $(echo $shell))"
